# extra/udev_rules/99-dwmblocks.rules
ACTION=="change", SUBSYSTEM=="power_supply", ATTR{type}=="Mains", RUN+="/bin/pkill -x -RTMIN+9 dwmblocks"
